# plugin/search/server/cron.d/dokeos-indexer
0 5 * * * root /usr/sbin/indexer -N 2 /etc/mnogosearch/indexer.conf &>/dev/null
